# systemd/bitcoind.service
[Unit]
Description=Bitcoin Daemon
After=network.target

[Service]
ExecStart=/bin/bitcoind

[Install]
WantedBy=default.target
