# systemd/canto.service
[Unit]
Description=Canto RSS Daemon
After=network.target

[Service]
ExecStart=/bin/canto-daemon

[Install]
WantedBy=default.target
