# systemd/ibus.service
[Unit]
Description=Ibus Daemon

[Service]
ExecStart=/usr/bin/ibus-daemon -x -r
Restart=always

[Install]
WantedBy=default.target
