# systemd/offlineimap.timer
[Unit]
Description=Run offlineimap full sync

[Timer]
OnBootSec=15
OnUnitActiveSec=5m

[Install]
WantedBy=default.target
