# systemd/rustdocs.timer
[Unit]
Description=Build rust doc epubs and push to git

[Timer]
OnCalendar=daily
Persistent=True

[Install]
WantedBy=default.target
