# systemd/rustdocs.service
[Unit]
Description=Build rust docs

[Service]
ExecStart=/home/chris/dev/bin/build_rust_docs.fish
